# etc/udev/rules.d/90-usb-libvirt-hotplug-psvita.rules
SUBSYSTEM=="usb",DEVPATH=="/devices/pci0000:00/0000:00:01.2/0000:02:00.0/0000:03:08.0/0000:09:00.3/usb5/5-1",RUN+="/opt/usb-libvirt-hotplug.sh"
